# systemd-units/syncthing.service.d/override.conf
# Keep Syncthing's config + device keys in ~/DATA/IT (synced, backed up),
# while the database/logs stay on local fast storage. See setup.sh "Syncthing".
[Service]
Environment="STCONFDIR=%h/DATA/IT/configs/syncthing"
Environment="STDATADIR=%h/.local/state/syncthing"
# Fedora ships syncthing 2.0.x (config schema v51); the synced config is v52,
# written by a 2.1.x device. The delta is benign (no v2.1-only elements present),
# so allow the newer config until Fedora ships 2.1.x.
ExecStart=
ExecStart=/usr/bin/syncthing serve --no-browser --no-restart --allow-newer-config
