# systemd-units/syncthing.service.d/override.conf
# Keep Syncthing's config + device keys in ~/DATA/IT (synced, backed up),
# while the database/logs stay on local fast storage. See setup.sh "Syncthing".
[Service]
Environment="STCONFDIR=%h/DATA/IT/configs/syncthing"
Environment="STDATADIR=%h/.local/state/syncthing"
