# systemd-units/sway-battery-warn.timer
[Unit]
Description=Sway battery check

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target

